# roles/fastapi/templates/fastapi.service.j2
[Unit]
Description=FastAPI Application {{ item.name }}
After=network.target

[Service]
User={{ deploy_user }}
Group={{ deploy_user }}
WorkingDirectory=/var/www/{{ item.name }}
ExecStart=/var/www/{{ item.name }}/venv/bin/uvicorn {{ item.module }}:{{ item.callable }} \
  --host 0.0.0.0 --port {{ item.port }} --workers {{ item.workers | default(2) }}
Restart=always

[Install]
WantedBy=multi-user.target
